# ztp/extra-manifests-builder/01-disk-encryption-pcr-rebind/pcr-disable-shutdown.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/true
ExecStop=/usr/local/bin/disablePcrOnRebootOrShutdown.sh
[Install]
WantedBy=multi-user.target
