# ztp/source-crs-builder/05-chrony-dynamic/chronyd-restart.timer
[Unit]
Description=Restart chronyd periodically

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
