# ztp/extra-manifests-builder/01-disk-encrypt-pcr/pcr-disable-shutdown.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/usr/local/bin/disablePcrOnRebootOrShutdown.sh

[Install]
WantedBy=multi-user.target
